# debian/prerm
#!/bin/sh

systemctl disable json-dbus-bridge@system.socket